# supervisor_conf/detection.conf
[program:detection_person]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map person"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_person.log
stdout_logfile_maxbytes=10MB

[program:detection_hand]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map hand"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_hand.log
stdout_logfile_maxbytes=10MB

[program:detection_fire]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map fire"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_fire.log
stdout_logfile_maxbytes=10MB

[program:detection_nest]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map nest"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_nest.log
stdout_logfile_maxbytes=10MB

[program:detection_kite]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map kite"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_kite.log
stdout_logfile_maxbytes=10MB

[program:detection_oil]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map oil"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_oil.log
stdout_logfile_maxbytes=10MB

[program:detection_smoke]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map smoke"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_smoke.log
stdout_logfile_maxbytes=10MB

[program:detection_helmet_suit_smoking]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map helmet_suit_smoking"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_helmet_suit_smoking.log
stdout_logfile_maxbytes=10MB

[program:detection_insulator]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map insulator"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_insulator.log
stdout_logfile_maxbytes=10MB

[program:detection_capacitor_bulge]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map capacitor_bulge"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_capacitor_bulge.log
stdout_logfile_maxbytes=10MB

[program:detection_box_door]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map box_door"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_box_door.log
stdout_logfile_maxbytes=10MB

[program:detection_blurred_dial]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map blurred_dial"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_blurred_dial.log
stdout_logfile_maxbytes=10MB

[program:detection_abnormal_meter]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map abnormal_meter"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_abnormal_meter.log
stdout_logfile_maxbytes=10MB

[program:detection_silicagel]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map silicagel"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_silicagel.log
stdout_logfile_maxbytes=10MB

[program:detection_screen_crash]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map screen_crash"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_screen_crash.log
stdout_logfile_maxbytes=10MB

[program:detection_damaged]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map damaged"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_damaged.log
stdout_logfile_maxbytes=10MB

[program:detection_oil_level]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map oil_level"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_oil_level.log
stdout_logfile_maxbytes=10MB

[program:detection_switch_1]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map switch_1"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_switch_1.log
stdout_logfile_maxbytes=10MB

[program:detection_switch_2]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map switch_2"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_switch_2.log
stdout_logfile_maxbytes=10MB


[program:detection_surface_pollution]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map surface_pollution"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_surface_pollution.log
stdout_logfile_maxbytes=10MB

[program:detection_rust]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map rust"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_rust.log
stdout_logfile_maxbytes=10MB

[program:detection_light_status]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map light_status"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_light_status.log
stdout_logfile_maxbytes=10MB

[program:detection_yaban]
command=/bin/bash -c "/home/dl/miniconda3/bin/python detection.py --model_type_map yaban"; 被监控的进程启动命令
directory=/home/dl/sh_cv_robot/                ; 执行前要不要先cd到目录去，一般不用
priority=1                    ;数字越高，优先级越高
numprocs=1                    ; 启动几个进程
autostart=true                ; 随着supervisord的启动而启动
autorestart=true              ; 自动重启。。当然要选上了
startretries=5             ; 启动失败时的最多重试次数
exitcodes=0                   ; 正常退出代码（是说退出代码是这个时就不再重启了吗？待确定）
stopasgroup=true
stopsignal=QUIT               ; 用来杀死进程的信号
stopwaitsecs=0              ; 发送SIGKILL前的等待时间
redirect_stderr=true          ; 重定向stderr到stdout
stdout_logfile=/tmp/detection_yaban.log
stdout_logfile_maxbytes=10MB